# .config/autostart/alacritty.desktop
[Desktop Entry]
Exec=/usr/bin/alacritty
Name=alacritty
Comment=A fast, cross-platform, OpenGL terminal emulator
Type=Application
