# .config/autostart/emacs.desktop
[Desktop Entry]
Exec=/usr/bin/emacs
Name=emacs
Comment=GNU project Emacs editor
Type=Application
