# .config/autostart/firefox.desktop
[Desktop Entry]
Exec=/usr/bin/firefox
Name=firefox
Type=Application
