# .config/autostart/mailspring.desktop
[Desktop Entry]
Exec=/usr/bin/mailspring
Name=mailspring
Type=Application
